# config/restic.service
[Unit]
Description=Restic backup service

[Service]
Type=oneshot
ExecStart=/home/ale/.local/bin/backup

[Install]
WantedBy=default.target